# geneva/geneva_https.path
[Unit]
Description=Control start stop of geneva https

[Path]
PathExists=/www/server/mdserver-web/data/geneva_https.pl

[Install]
WantedBy=multi-user.target
