# geneva/geneva_http.path
[Unit]
Description=Control start stop of geneva https

[Path]
PathExists=/www/server/mdserver-web/data/geneva_http.pl

[Install]
WantedBy=multi-user.target
