# scripts/speed/speed_check.service
[Unit]
Description=speed_check

[Service]
ExecStart=/www/server/mdserver-web/scripts/speed/speed_check.sh